# strange-lock.service
# ~/.config/systemd/user/strange-screensaver.service
[Unit]
Description=Strange attractor screensavers (glslViewer)
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/share/screensaver/run-simone.sh
ExecStop=/usr/bin/pkill -f 'glslViewer.*simone.frag'
Restart=on-failure

[Install]
WantedBy=default.target
